# ntp-data-collector.service
[Unit]
Description=NTP Data Collection

[Service]
Type=simple
ExecStart=/lib/rdk/ntp-data-collector.sh

[Install]
WantedBy=multi-user.target
